# connman/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/dbus-1/
drwxr-xr-x	root/root	etc/dbus-1/system.d/
-rw-r--r--	root/root	etc/dbus-1/system.d/connman.conf
drwxr-xr-x	root/root	etc/systemd/
drwxr-xr-x	root/root	etc/systemd/system/
-rw-r--r--	root/root	etc/systemd/system/connman.service
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/connman/
-rw-r--r--	root/root	usr/include/connman/device.h
-rw-r--r--	root/root	usr/include/connman/inet.h
-rw-r--r--	root/root	usr/include/connman/ipconfig.h
-rw-r--r--	root/root	usr/include/connman/log.h
-rw-r--r--	root/root	usr/include/connman/network.h
-rw-r--r--	root/root	usr/include/connman/notifier.h
-rw-r--r--	root/root	usr/include/connman/plugin.h
-rw-r--r--	root/root	usr/include/connman/provision.h
-rw-r--r--	root/root	usr/include/connman/resolver.h
-rw-r--r--	root/root	usr/include/connman/service.h
-rw-r--r--	root/root	usr/include/connman/storage.h
-rw-r--r--	root/root	usr/include/connman/types.h
-rw-r--r--	root/root	usr/include/connman/version.h
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/connman/
drwxr-xr-x	root/root	usr/lib/connman/scripts/
drwxr-xr-x	root/root	usr/lib/connman/test/
-rwxr-xr-x	root/root	usr/lib/connman/test/backtrace
-rwxr-xr-x	root/root	usr/lib/connman/test/connect-vpn
-rwxr-xr-x	root/root	usr/lib/connman/test/disable-tethering
-rwxr-xr-x	root/root	usr/lib/connman/test/disconnect-vpn
-rwxr-xr-x	root/root	usr/lib/connman/test/enable-tethering
-rwxr-xr-x	root/root	usr/lib/connman/test/get-global-timeservers
-rwxr-xr-x	root/root	usr/lib/connman/test/get-proxy-autoconfig
-rwxr-xr-x	root/root	usr/lib/connman/test/get-services
-rwxr-xr-x	root/root	usr/lib/connman/test/get-state
-rwxr-xr-x	root/root	usr/lib/connman/test/list-services
-rwxr-xr-x	root/root	usr/lib/connman/test/monitor-connman
-rwxr-xr-x	root/root	usr/lib/connman/test/monitor-services
-rwxr-xr-x	root/root	usr/lib/connman/test/service-move-before
-rwxr-xr-x	root/root	usr/lib/connman/test/set-domains
-rwxr-xr-x	root/root	usr/lib/connman/test/set-global-timeservers
-rwxr-xr-x	root/root	usr/lib/connman/test/set-ipv4-method
-rwxr-xr-x	root/root	usr/lib/connman/test/set-ipv6-method
-rwxr-xr-x	root/root	usr/lib/connman/test/set-nameservers
-rwxr-xr-x	root/root	usr/lib/connman/test/set-proxy
-rwxr-xr-x	root/root	usr/lib/connman/test/show-introspection
-rwxr-xr-x	root/root	usr/lib/connman/test/simple-agent
-rwxr-xr-x	root/root	usr/lib/connman/test/test-clock
-rwxr-xr-x	root/root	usr/lib/connman/test/test-compat
-rwxr-xr-x	root/root	usr/lib/connman/test/test-connman
-rwxr-xr-x	root/root	usr/lib/connman/test/test-counter
-rwxr-xr-x	root/root	usr/lib/connman/test/test-manager
-rwxr-xr-x	root/root	usr/lib/connman/test/test-new-supplicant
-rwxr-xr-x	root/root	usr/lib/connman/test/test-session
-rwxr-xr-x	root/root	usr/lib/connman/test/test-supplicant
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/connman.pc
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/connmand
